# SGMsensor/SGMsensor.service/SGMsensor.service
[Unit]
Description = SGM Sensor service

[Service]
Type          = simple
ExecStart     = /home/root/sgm/SGMsensor
